# server_configurations/storage/etc/service/storage-synchronization/run
#!/bin/bash

# start a storage writer
cd /home/indextank/storage
echo "Supervise's run starting synchronization"

while [ 1 ]; do
  sudo -u indextank /bin/bash ./synchronize.sh 2>&1
  sleep 30
done
